# ls_backup/scripts/autostart/sikuli_setup.sh.desktop
[Desktop Entry]
Type=Application
Exec=/home/vagrant/linux_setup/sikuli/sikuli_setup.sh
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name[en_US]=Sikuli Setup
Name=Sikuli Setup
Comment[en_US]=
Comment=
